# speedflux/rootfs/etc/services.d/example/run
#!/usr/bin/env bashio
# shellcheck shell=bash
set +u
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

## Get the 'message' key from the user config options.


bashio::log.info "Starting /app/main.py..."

if bashio::config.has_value 'namespace'; then
    namespace=$(bashio::config 'namespace')
    export NAMESPACE=${namespace}
fi

if bashio::config.has_value 'influx_db_address'; then
    influx_db_address=$(bashio::config 'influx_db_address')
    export INFLUX_DB_ADDRESS=${influx_db_address}
fi

if bashio::config.has_value 'influx_db_port'; then
    influx_db_port=$(bashio::config 'influx_db_port')
    export INFLUX_DB_PORT=${influx_db_port}
fi

if bashio::config.has_value 'influx_db_user'; then
    influx_db_user=$(bashio::config 'influx_db_port')
    export INFLUX_DB_USER=${influx_db_user}
fi

if bashio::config.has_value 'influx_db_password'; then
    influx_db_password=$(bashio::config 'influx_db_password')
    export INFLUX_DB_PASSWORD=${influx_db_password}
fi

if bashio::config.has_value 'influx_db_database'; then
    influx_db_database=$(bashio::config 'influx_db_database')
    export INFLUX_DB_DATABASE=${influx_db_database}
fi

if bashio::config.has_value 'influx_db_tags'; then
    influx_db_tags=$(bashio::config 'influx_db_tags')
    export INFLUX_DB_TAGS=${influx_db_tags}
fi

if bashio::config.has_value 'speedtest_interval'; then
    speedtest_interval=$(bashio::config 'speedtest_interval')
    export SPEEDTEST_INTERVAL=${speedtest_interval}
fi

if bashio::config.has_value 'speedtest_server_id'; then
    speedtest_server_id=$(bashio::config 'speedtest_server_id')
    export SPEEDTEST_SERVER_ID=${speedtest_server_id}
fi

if bashio::config.has_value 'ping_interval'; then
    ping_interval=$(bashio::config 'ping_interval')
    export PING_INTERVAL=${ping_interval}
fi

if bashio::config.has_value 'ping_targets'; then
    ping_targets=$(bashio::config 'ping_targets')
    export PING_TARGETS=${ping_targets}
fi

if bashio::config.has_value 'tweaks.log_type'; then
    bashio::log.info "log_type ha value!"
    log_type=$(bashio::config 'tweaks.log_type')
    export LOG_TYPE=${log_type}
fi

bashio::log.info "log_type: ${LOG_TYPE}"
bashio::log.info "ping_targets: ${PING_TARGETS}"

## Run your program
exec python3 /app/main.py
